# config/systemd/user/profile-sync.target
##
## ~/.config/systemd/user/profile-sync.target
##
[Unit]
Description=Profile Synchronization Target
StopWhenUnneeded=yes
